# postgrey/supervisord.ini
[supervisord]
nodaemon=true

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:postgrey]
command=/usr/sbin/postgrey --inet=127.0.0.1:10023 --pidfile=/var/run/postgrey.pid --delay=50 --greylist-text="Greylisted for %s seconds"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
